# battery_status_led.service
[Unit]
Description=Battery monitoring on CapsLock LED

[Service]
ExecStart="/usr/bin/battery_led.sh"
Type=idle

